# task.cron
*/12 * * * * /opt/backup-bitnami-redmine.sh
