# task.cron
0 6 * * * /opt/bitnami-redmine-backup/backup.sh
